# setup/changeme.service
[Unit]
Description=changeme

[Service]
ExecStart=/usr/bin/npm run deploy --prefix /srv/changeme/
Restart=always

[Install]
WantedBy=multi-user.target